# supervisord.conf
[supervisord]
loglevel=debug
logfile=/dev/stdout
logfile_maxbytes=0

[inet_http_server]
port=127.0.0.1:9001

# proxy components inside enclave (host CID=3)
[program:ip-to-vsock-raw-outgoing]
command=/enclaved/ip-to-vsock-raw-outgoing --vsock-addr 3:1080 --queue-num 0
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# enclave CID must be 16
[program:vsock-to-ip-raw-incoming]
command=/enclaved/vsock-to-ip-raw-incoming --vsock-addr 16:1080 --device br0
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# DNS-over-HTTPS provider
[program:dnsproxy]
command=/enclaved/dnsproxy -u https://1.1.1.1/dns-query
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# docker
[program:docker]
command=/bin/dockerd
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# phoenixd
[program:phoenixd]
command=/enclaved/enclave-phoenixd.sh
directory=/home/phoenix
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0

# enclaved
[program:enclaved]
command=/enclaved/enclave-enclaved.sh
autostart=false
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stdout
stderr_logfile_maxbytes=0